# dotfile-etc.list
# SYNTAX: COMMAND <item address without profile folder address>
# COMMAND:
#     l: item to be linked
#     c: item to be copied
# NOTE: - ITEM address should without spaces on name
#       - comments should be on their own lines
#       - some items are still not on the repo yet

c sysctl.d
c makepkg.conf
c udev/rules.d
c modules-load.d
c systemd/system/user@.service.d/dbus.conf
c systemd/system/sshdgenkeys.service
c systemd/system/tmux@.service
c systemd/system/dnsmasq.service
c systemd/system/docker.service
c systemd/system/connman-vpn.service
c systemd/journald.conf
c systemd/logind.conf
c systemd/user/dbus.socket
c systemd/user/dbus.service
c abs.conf
c default/grub
c locale.conf
c mkinitcpio.d
c psd.conf
c asd.conf
c pam.d/su
c pam.d/system-login
c initcpio
c modprobe.d
c mkinitcpio-lts.conf
c mkinitcpio.conf
c fstab
c yaourtrc
c customizepkg.d
#c davfs2/davfs2.conf
c bitlbee
#c polipo
c pacman.conf
c pacman.d/mirrorlist
c vconsole.conf
c localtime
c hostname
c connman
c systemd/network
c vimrc
c X11/xorg.conf.d
c slim.conf
c slimlock.conf
c dnsmasq.d
c dnsmasq.conf
c dhcpcd.conf
c resolv.conf
c hosts
c hostapd
c hosts.ap0
c issue
c qemu-ifup
c qemu-ifdown
c hosts.qemu0
c nodogsplash
c ssh/ssh_config
c ssh/sshd_config

# wpa_supplicant
# profile
# skel
# profile.d
# host.conf
# ssl/openssl.cnf
# openldap
# security
# login.defs
