# ~/.config/systemd/user/getmails.timer
[Unit]
Description=2 Min updates for user mails

[Timer]
OnBootSec=1min
OnUnitActiveSec=50min
Persistent=true

[Install]
WantedBy=default.target
