# ~/.config/systemd/user/modprobed-db.timer
[Unit]
Description=10 min check for probed modules

[Timer]
OnBootSec=1min
OnUnitActiveSec=10min
Persistent=true

[Install]
WantedBy=default.target
