# etc/udev/rules.d/70-disable_wol.rules
ACTION=="add", SUBSYSTEM=="net", KERNEL=="enp9s*", RUN+="/usr/bin/ethtool -s %k wol d"
