# etc/udev/rules.d/70-wifi-powersave.rules
ACTION=="add", SUBSYSTEM=="net", KERNEL=="wlp8s*", RUN+="/usr/bin/iw dev %k set power_save on"
