# etc/systemd/system/virtualbox.service
[Unit]
Description=Virtualbox setup

[Service]
ExecStart=/usr/bin/rcvboxdrv setup

[Install]
WantedBy=multi-user.target
